# modules/ganglia/files/gmetad.chkconfig
# chkconfig: 2345 98 02
# description: gmetad startup script
