# modules/varnish/files/logrotate-vb53u.varnish.cron
#!/bin/sh

# $Id$
# $URL$

/usr/bin/logger -t logrotate "Starting log varnish rotation"
/usr/sbin/logrotate /etc/logrotate.d/varnish
EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
    /usr/bin/logger -t logrotate "ALERT exited abnormally with [$EXITVALUE]"
fi
/usr/bin/logger -t logrotate "log rotation exited normally"
exit 0
