# modules/cronjob/files/cron.d/crowdignite_engine/eng6.cron
# those in .sh scripts need to be run in serial

# These jobs read from the new widget_stats_views and landing_page_stats_views queues
*/2 * * * * root sudo -u apache /usr/bin/php -c /etc/php_browscap_cli.ini -q /app/shared/public_html/app/cron_dispatcher.php /stats_engine/update_widget_stats_impressions/4/20000/100/1 >> /app/log/cronjobs/widget_update_stats_impressions_4.log 2>&1
*/2 * * * * root sudo -u apache /usr/bin/php -c /etc/php_browscap_cli.ini -q /app/shared/public_html/app/cron_dispatcher.php /stats_engine/update_widget_stats_impressions/8/20000/100/1 >> /app/log/cronjobs/widget_update_stats_impressions_8.log 2>&1
*/5 * * * * root sudo -u apache /usr/bin/php -q /app/shared/public_html/app/cron_dispatcher.php /stats_engine/update_landing_page_widget_tracking/1/1000/100/1 >> /app/log/cronjobs/update_landing_page_widget_tracking_1.log 2>&1

#
*/5 * * * * root sudo -u apache /usr/bin/php -q /app/shared/public_html/app/cron_dispatcher.php /landing_page_engine/update_stats_impressions/4 >> /app/log/cronjobs/lp_update_stats_impressions_4.log 2>&1
*/5 * * * * root sudo -u apache /usr/bin/php -q /app/shared/public_html/app/cron_dispatcher.php /landing_page_engine/update_stats_impressions/8 >> /app/log/cronjobs/lp_update_stats_impressions_8.log 2>&1
