# modules/cronjob/files/cron.d/test/test.cron
* * * * * root /usr/local/bin/test/test.sh
