# modules/cronjob/files/cron.d/sphinx_vb/sphinx.cron
# $Id$
# $URL$

# make sure full index is not overlab
*/5 * * * * root /usr/local/bin/sphinx_vb/sphinx-delta.sh 
0 3 * * * root /usr/local/bin/sphinx_vb/sphinx-full-index.sh 
