# modules/cronjob/files/cron.d/vb/tfs_vbseo_sitemap.cron
# $Id$
# $URL$

10 2 * * * root php /app/shared/forums.thefashionspot.com/vbseo_sitemap/vbseo_sitemap.php 2>&1
